# Languages/Python/Frame/Supervisor/Storages/study3.conf
[program:study3]
user = ubuntu
command=/home/ubuntu/Code/Python/virtual/py3env/bin/uwsgi --ini uwsgi.ini
directory = /home/ubuntu/Github/python_study/django/study3
startsecs = 0
stopwaitsecs = 0
autorstart= true
autorestart = true
redirect_stderr = true
stdout_logfile = /home/ubuntu/Github/python_study/django/study3/logs/study3.super.log
stderr_logfile = /home/ubuntu/Github/python_study/django/study3/logs/study3.super.err
[supervisord]
loglevel = info
